# supervisord.conf
[supervisord]
nodaemon=true

[program:ssr]
command=node /app/dist/server/main.js
autostart=true
autorestart=true

[program:goapp]
command=/app/jgo
autostart=true
autorestart=true